# config/etc/systemd/system/crash-reporter@.service
[Unit]
Description=Meticulous Crash Reporter Service
After=network.target

[Service]
Type=oneshot
ExecStart=/etc/crash-reporter/crash-reporter.sh