# system-services/meticulous-watcher.service
[Unit]
Description=Meticulous Watcher
After=network.target
Before=weston.service

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/python3.12 /opt/meticulous-watcher.py
WantedBy=multi-user.target
User=root
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target
