# system-services/meticulous-smoke-report.service
[Unit]
Description=Publish Meticulous smoke validation status to Hawkbit
ConditionHost=meticulous*
After=network-online.target meticulous-backend.service meticulous-dial.service rauc-hawkbit-updater.service
Wants=network-online.target
OnFailure=crash-reporter@.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/meticulous-smoke-report --boot
TimeoutStartSec=4min
User=root

[Install]
WantedBy=multi-user.target
